# cache_nodes_012426_2236/configs/systemd/node-exporter.service
[Unit]
Description=Prometheus Node Exporter
After=network.target

[Service]
User=nobody
ExecStart=/usr/local/bin/node_exporter \
  --collector.disable-defaults \
  --collector.cpu \
  --collector.meminfo \
  --collector.filesystem \
  --collector.netdev \
  --collector.loadavg \
  --collector.time \
  --collector.systemd \
  --collector.filesystem.mount-points-exclude="^/(proc|sys|dev|tmp|run|var/log)($|/)" \
  --web.listen-address=192.168.86.20:9100
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target

# Notes:
# - Change --web.listen-address to 192.168.86.19:9100 for backend node
# - Collectors are minimal to reduce CPU/memory usage
# - Excludes tmpfs mounts from filesystem metrics
