# cache_nodes_012426_2236/configs/systemd/node-exporter.service
[Unit]
Description=Prometheus node exporter (stub)
After=network-online.target

[Service]
ExecStart=/usr/bin/node_exporter
# RAM-safe options; deploy per Phase 5.
