# cache_nodes_012426_2236/configs/systemd/okome-health-probe.service
[Unit]
Description=OKOME upstream health probe (writes nginx include)
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/okome-health-probe.sh
User=root
