# cache_nodes_012426_2236/configs/systemd/okome-health-probe.service
# OKOME – health probe (stub)
# Phase 2. Unit for okome-health-probe.timer.

[Unit]
Description=OKOME health probe
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/okome-health-probe.sh
