# cache_nodes_012426_2236/configs/systemd/okome-health-probe.timer
[Unit]
Description=Run OKOME health probe every 5 seconds

[Timer]
OnBootSec=5s
OnUnitActiveSec=5s
AccuracySec=1s

[Install]
WantedBy=timers.target
