# cache_nodes_012426_2236/configs/systemd/okome-prewarm.timer
[Unit]
Description=OKOME pre-warm timer (stub)
[Timer]
OnBootSec=60
OnUnitActiveSec=120
Persistent=yes
