# cache_nodes_012426_2236/configs/systemd/okome-prewarm.service
[Unit]
Description=OKOME pre-warm (stub)
After=network-online.target nginx.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/okome-prewarm.sh
